# deploy/vframotp-worker.service
[Unit]
Description=VFRAMOTP Worker
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=vfarmotp
WorkingDirectory=/opt/vframotp
EnvironmentFile=/opt/vframotp/.env
ExecStart=/opt/vframotp/bin/worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
